# data/pam_biod.service
[Unit]
Description=PamBio Daemon

[Service]
Type=dbus
BusName=xyz.ccat3z.pambio
ExecStart=@PAM_BIOD_PATH@
